# services/systemd/homelab-backup.service
[Unit]
Description=Homelab Unified Backup Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=egk
Group=egk
Environment=HOME=/home/egk
WorkingDirectory=/home/egk/homelab
ExecStart=/home/egk/homelab/scripts/backup/unified_backup.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
